# openwrt/target/linux/ramips/base-files/etc/init.d/smp
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=99
start() {
	echo "run smp.sh...."
	/sbin/smp.sh
}

