# openwrt/package/base-files/files/etc/crontabs/root
* * * * * /etc/init.d/wifi_control start
